# sendspin/rootfs/etc/services.d/sendspin/run
#!/usr/bin/with-contenv bashio

declare -a options
declare player_name
declare player_id
declare static_delay_ms

bashio::log.info 'Starting sendspin daemon...'

# Player Name
bashio::config.require 'player_name'
player_name="$(bashio::config 'player_name')"
options+=(--name "$player_name")

# Player ID
if bashio::config.has_value 'player_id'; then
    player_id="$(bashio::config 'player_id')"
    bashio::log.debug "Using configured player id: $player_id"
    options+=(--id "$player_id")
fi

# Static delay
if bashio::config.has_value 'static_delay_ms'; then
    static_delay_ms="$(bashio::config 'static_delay_ms')"
    bashio::log.debug "Setting extra delay: ${static_delay_ms}ms"
    options+=(--static-delay-ms "$static_delay_ms")
fi

# Debug Mode
if bashio::debug; then
    options+=(--log-level DEBUG)
fi

options+=(--disable-mpris)
options+=(--manufacturer "Adrian Freund")
options+=(--product-name "Home Assistant Sendspin App")

exec /venv/bin/sendspin daemon "${options[@]}"
